# api.service
[Unit]
Description=

[Service]
User=
Group=
WorkingDirectory=
ExecStart=
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
